# contrib/systemd/wob.socket
[Socket]
ListenFIFO=%t/wob.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
